# deploy/yao-gpt-frontend.service
[Unit]
Description=Yao GPT Gradio Frontend
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h/repos/yao-gpt-service
Environment="PATH=%h/repos/yao-gpt-service/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/repos/yao-gpt-service/.venv/bin/python frontend/run_frontend.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
